# www/tuxedocat-redirect.socket
[Unit]
Description=TuxedoCat Lounge BBS redirect message

[Socket]
ListenStream=2323
Accept=yes

[Install]
WantedBy=sockets.target
